# templates/upstart-clark-kent.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          clark-kent
# Required-Start:    $local_fs $network $syslog $named
# Required-Stop:     $local_fs $network $syslog $named
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# X-Interactive:     true
# Short-Description: Start/stop clark-kent ILS reporter
# Description:       Start the reporter deamon
#  This script will start and stop clark-kent.p .
### END INIT INFO

DESC="clark-kent reporter deamon"
NAME=clark-kent.pl
DAEMON=<%= sysdir %>/bin/$NAME

USER=<%= user %>
SLEEP=<%= clark_kent_sleep %>
LOCKFILE=<%= clark_kent_lockfile %>
CONCURRENCY=<%= clark_kent_concurrency %>
BOOTSTRAP=<%= sysdir %>/conf/opensrf_core.xml

# sleep=interval: number of seconds to sleep between checks for new reports to run; defaults to 10
# lockfile=filename: where to place the lockfile for the process; defaults to /tmp/reporter-LOCK
# concurrency=integer: number of reporter daemon processes to run; defaults to 1
# bootstrap=filename: OpenSRF bootstrap

start() {
    echo "Starting $DESC" "$NAME"
    sudo -u $USER $DAEMON --daemon -sleep=$SLEEP -lockfile=$LOCKFILE -concurrency=$CONCURRENCY -bootstrap=$BOOTSTRAP
}

stop() {
    echo "Stopping $DESC" "$NAME"
    r=$(ps wax | grep "Clark Kent" | grep -v grep | cut -b1-6)
    if [ ! -z "$r" ] ; then
        kill $r
    fi
    if [ -f $LOCKFILE ] ; then
        rm $LOCKFILE
    fi
}

restart() {
    stop;
    start;
}

case "$1" in
  start)
	start;
	;;
  stop)
    stop;
    ;;
  restart)
    restart;
    ;;
  *)
	echo "Usage: service clark-kent start|stop|restart" >&2
	exit 3
	;;
esac



exit 0